# wsl-docker/setContainer/docker/supervisord/postfix.ini
[group:postfix]
programs=postfixd
priority=30

[program:postfixd]
command=/usr/sbin/postfix start-fg
process_name=%(program_name)s
startsecs = 0
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
